# schedule.cron
#SPRINKLER # Garden Box Bubblers on Monday
0 5 * * 1 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=1&time=1"
#SPRINKLER # Garden Box Bubblers on Thursday
0 5 * * 4 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=1&time=1"
#SPRINKLER # Left Rear Garden Box on Monday
0 5 * * 1 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=2&time=30"
#SPRINKLER # Left Rear Garden Box on Thursday
0 5 * * 4 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=2&time=30"
#SPRINKLER # Rear Right Veg and Berm on Thursday
0 5 * * 4 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=4&time=30"
#SPRINKLER # Rear Rose Garden on Monday
0 5 * * 1 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=8&time=30"
#SPRINKLER # Rear Rose Garden on Thursday
0 5 * * 4 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=8&time=30"
#SPRINKLER # Left Redwood Sprays on Thursday
0 5 * * 4 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=9&time=30"
#SPRINKLER # Rear Right Redwoods (blue) on Thursday
0 5 * * 4 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=12&time=30"
#SPRINKLER # Rear Near Spray on Monday
0 5 * * 1 /usr/bin/curl "http://10.0.0.241:8080/api/enqueue?id=7&time=15"
